# configs/com.macos.guest-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macos.guest-agent</string>
    
    <key>Program</key>
    <string>/usr/local/bin/mac-guest-agent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/mac-guest-agent</string>
        <string>--daemon</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>UserName</key>
    <string>root</string>
    
    <key>GroupName</key>
    <string>wheel</string>
    
    <key>StandardOutPath</key>
    <string>/var/log/mac-guest-agent.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/mac-guest-agent.log</string>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/share/mac-guest-agent</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist> 